# systemd/.config/systemd/user/backup.service
[Unit]
Description=Daily rsync backup (user) 

[Service]
Type=oneshot
ExecStart=%h/.local/bin/backup.sh
